# airprint/rootfs/etc/s6-overlay/s6-rc.d/airprint-updater/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2034

ulimit -n 524288
bashio::log.info "Starting Avahi daemon from S6"

debug_flag="false"
if bashio::config.has_value 'avahi_debug'; then
  debug_flag=$(bashio::config 'avahi_debug')
fi

# shellcheck disable=SC2034
avahi_services=/config/avahi/services
cups_config=/config/cups

/usr/bin/inotifywait -m -e close_write,moved_to,create "$cups_config" |
  while read -r directory events filename; do
    if [ "$filename" = "printers.conf" ]; then
      bashio::log.info "Running AirPrint update"
      rm -f "$avahi_services"/AirPrint-*.service

      if [ "$debug_flag" == true ]; then
        /opt/airprint/airprint-generate.py --verbose --host=127.0.0.1 --port=631 --directory="$avahi_services" --prefix=AirPrint-
      else
        /opt/airprint/airprint-generate.py --host=127.0.0.1 --port=631 --directory="$avahi_services" --prefix=AirPrint-
      fi
    fi
  done
