# rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288
bashio::log.info "Starting Avahi daemon from S6"

debug_flag="false"
if bashio::config.has_value 'avahi_debug'; then
  debug_flag=$(bashio::config 'avahi_debug')
fi

until [ -e /run/dbus/system_bus_socket ]; do
  bashio::log.info "Waiting for dbus daemon."
  sleep 2s
done

if [ "$debug_flag" == "true" ]; then
  avahi-daemon --debug -f /config/avahi/avahi-daemon.conf 
else
  avahi-daemon --file=/config/avahi/avahi-daemon.conf
fi

