# services/cups-server/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2034

# shellcheck source="../../src/common/paths/cups-paths.sh"
source "/opt/common/paths/cups-paths.sh"

ulimit -n 524288

function run() {
    /usr/bin/inotifywait -m -e close_write,moved_to,create "$real_cups_path" |
        while read -r directory events filename; do
            if [ "$filename" = "$cups_daemon" ] || [ "$filename" = "$cups_files" ] || [ "$filename" = "$cups_client" ]; then

                bashio::log.info "Testing CUPS server config"
                cupsd -t -c "$real_cups_path"/"$cups_daemon" -s "$real_cups_path"/"$cups_files"

                bashio::log.info "Starting CUPS server from S6"

                cupsd -l "$real_cups_path"/"$cups_daemon" -s "$real_cups_path"/"$cups_files"
            fi
        done
}

run
