# services/cups-server/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2034

# shellcheck source="../../src/common/paths.sh"
source "/opt/common/paths.sh"

ulimit -n 524288

function check_install() {
    if command -v cupsd &>/dev/null; then
        bashio::log.debug "Cupsd is installed"
    else
        bashio::log.debug "Installing Cups Libs"
        find /build -type f -name "cups-libs-$CUPS_VER-linux-**.deb" -exec bash -c 'for pkg; do dpkg -i --force-confold --force-confdef "$pkg" ; done' _ {} +

        bashio::log.debug "Installing Cups"
        find /build -type f -name "cups-$CUPS_VER-linux-**.deb" -exec bash -c 'for pkg; do dpkg -i --force-confold --force-confdef "$pkg"; done' _ {} +

        bashio::log.debug "Installing lib cups filters"
        apt-get update && apt-get install libcupsfilters2 -y --no-install-recommends
    fi
}

function run() {
    /usr/bin/inotifywait -m -e close_write,moved_to,create "$real_cups_path" |
        while read -r directory events filename; do
            if [ "$filename" = "$cups_daemon" ] || [ "$filename" = "$cups_files" ] || [ "$filename" = "$cups_client" ]; then

                bashio::log.info "Testing CUPS server config"
                cupsd -t -c "$real_cups_path"/"$cups_daemon" -s "$real_cups_path"/"$cups_files"

                bashio::log.info "Starting CUPS server from S6"
                cupsd -f -c "$real_cups_path"/"$cups_daemon" -s "$real_cups_path"/"$cups_files"
            fi
        done
}

check_install
run
