# services/initialization/run
#!/command/with-contend bashio

# bashio::log.info "Getting HA config"
# bashio /opt/helpers/ha-config-loader.sh

bashio::log.info "Initializing Common Config"
bashio /opt/common/common.sh

bashio::log.info "Initializing Avahi files"
bashio /opt/avahi/avahi-files.sh

bashio::log.info "Initializing Cups files"
bashio /opt/cups/cups-files.sh

bashio::log.info "Initializing NGINX files"
bashio /opt/nginx/nginx-files.sh

bashio::log.info "Running Additional packages"
bashio /opt/common/packages.sh

bashio::log.info "Updating Users"
bashio /opt/common/users.sh

bashio::log.info "Initializing Avahi config"
bashio /opt/avahi/avahi-config.sh

bashio::log.info "Initializing Cups config"
bashio /opt/cups/cups-config.sh

bashio::log.info "Initializing NGINX config"
bashio /opt/nginx/nginx-config.sh
