# services/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2034

# shellcheck source="../../src/common/paths.sh"
source "/opt/common/paths.sh"

ulimit -n 524288

function run() {
    nginx -g "daemon off;"
}

run
