# services/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2034

# shellcheck source="../../src/common/paths/nginx-paths.sh"
source "/opt/common/paths/nginx-paths.sh"

ulimit -n 524288

function run() {
    /usr/bin/inotifywait -m -e close_write,moved_to,create "$nginx_config_path" |
        while read -r directory events filename; do
            if [ "$filename" = "$nginx_default" ] || [ "$filename" = "$nginx_conf" ]; then
                bashio::log.info "Getting NGINX IPs"
                bashio /opt/helpers/nginx-ip-helper.sh

                bashio::log.info "Testing NGINX server config"
                nginx -t

                bashio::log.info "Starting NGINX server from S6"
                nginx -g "daemon off;"
            fi
        done
}

run
