# rootfs/etc/s6-overlay/s6-rc.d/systemd-resolved/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288
bashio::log.info "Starting systemd-resolved daemon from S6"

# Start systemd-resolved in the foreground
/usr/lib/systemd/systemd-resolved
