# rootfs/etc/s6-overlay/s6-rc.d/cups-browsed/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288
bashio::log.info "Starting cups-browsed daemon from S6"

/usr/sbin/cups-browsed
