# services/avahi-daemon/run
#!/usr/bin/with-contenv bashio
# shellcheck source="../../src/common/paths/avahi-paths.sh"
source "/opt/common/paths/avahi-paths.sh"

ulimit -n 524288

debug_flag="false"
if bashio::config.has_value 'AVAHI_SETTINGS.AVAHI_DEBUG'; then
    bashio::log.info "Setting Avahi debug flag"
    debug_flag=$(bashio::config 'AVAHI_SETTINGS.AVAHI_DEBUG')
else
    debug_flag="false"
fi

until [ -e  /run/dbus/system_bus_socket ]; do
    bashio::log.info "Waiting for dbus daemon before installing custom script"
    sleep 10s
done

bashio::log.info "Starting Avahi daemon from S6"
if [ "$debug_flag" == "true" ]; then
    avahi-daemon --debug --file="$avahi_config_path"/"$avahi_daemon"
else
    avahi-daemon --file="$avahi_config_path"/"$avahi_daemon"
fi
