# services/avahi-daemon/run
#!/usr/bin/with-contenv bashio
# shellcheck source="../../src/common/paths.sh"
source "/opt/common/paths.sh"

ulimit -n 524288

debug_flag="false"
if bashio::config.has_value 'avahi_debug'; then
    debug_flag=$(bashio::config 'avahi_debug')
fi

bashio::log.info "Starting Avahi daemon from S6"
if [ "$debug_flag" == "true" ]; then
    avahi-daemon --debug --file="$avahi_config_path"/"$avahi_daemon"
else
    avahi-daemon --file="$avahi_config_path"/"$avahi_daemon"
fi
